# scoped-uuplugin.service
[Unit]
Description=Scoped netease UU plugin.

[Service]
ExecStart=/usr/lib/scoped-uuplugin/daemon.sh

[Install]
WantedBy=multi-user.target
